# crontab
0 15 * * fstrim /
0 16 1 * vkpurge rm all
