# crontab
0 15 * * /sbin/fstrim /
0 16 1 * vkpurge rm all
