# etc/crontab
15    13     *      *       1                 /sbin/fstrim --all
